# examples/supervisord.conf
[program:pushycat]
command=/path/to/pushycat -conf /path/to/hooks.json
user=pushycat

autostart=true
autorestart=true
redirect_stderr=true
